# scripts/islamic-bot-update.service
[Unit]
Description=Weekly update for Islamic Unified Bot

[Service]
Type=oneshot
ExecStart=/opt/islamic-unified-bot/scripts/weekly-update.sh
User=root
