# scripts/islamic-bot-update.timer
[Unit]
Description=Weekly update timer for Islamic Unified Bot

[Timer]
OnCalendar=weekly
Persistent=true
RandomizedDelaySec=3600

[Install]
WantedBy=timers.target
